# docker/supervisor/laravel.conf
[program:laravel-serve]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan discord:serve
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
user=root
numprocs=1

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --sleep=3 --tries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
user=root
numprocs=5

[program:laravel-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/scheduler.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
user=root
numprocs=1
